# Adguard Home/adguard-log-downloader.timer
[Unit]
Description=Download query logs from AdGuard Home every minute

[Timer]
OnCalendar=*:0/1
Persistent=true
AccuracySec=10s

[Install]
WantedBy=timers.target
